# sw/99-sc4hub.rules
# Teknic SC4-Hub: use Exar xr_usb_serial driver instead of cdc_acm
#
# Install:
#   sudo mkdir -p /lib/modules/$(uname -r)/extra
#   sudo cp third_party/ExarKernelDriver/xr_usb_serial_common.ko /lib/modules/$(uname -r)/extra/
#   sudo depmod
#   sudo cp 99-sc4hub.rules /etc/udev/rules.d/
#   sudo udevadm control --reload-rules
#   # Then re-plug the SC4-Hub USB cable.

# When SC4-Hub appears, unbind cdc_acm and bind the Exar driver
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="2890", ATTR{idProduct}=="0213", \
    RUN+="/sbin/modprobe xr_usb_serial_common", \
    RUN+="/bin/sh -c 'sleep 0.5; for iface in /sys/bus/usb/drivers/cdc_acm/1-*; do echo $$(basename $$iface) > /sys/bus/usb/drivers/cdc_acm/unbind 2>/dev/null; done; echo %k:1.0 > /sys/bus/usb/drivers/cdc_xr_usb_serial/bind 2>/dev/null'"

# Give uucp group access and create symlink
SUBSYSTEM=="tty", ATTRS{idVendor}=="2890", ATTRS{idProduct}=="0213", MODE="0660", GROUP="uucp", SYMLINK+="sc4hub"
